# files/systemd/mirror@.service
[Unit]
Description=mirror %i

[Service]
Type=simple
ExecStart=/home/mirror/scripts/%i

[Install]
WantedBy=multi-user.target
