# templates/systemd/mirror@.service.erb
[Unit]
Description=mirror %i

[Service]
Type=simple
ExecStart=<%= @userhome %>/scripts/%i

[Install]
WantedBy=multi-user.target
