# templates/systemd/mirror@.service.erb
[Unit]
Description=mirror %i

[Service]
Type=simple
ExecStart=<%= @userhome %>/scripts/sync.sh %i
User=<%= @username %>
Group=<%= @usergroup %>

[Install]
WantedBy=multi-user.target
